# roles/bz_review_report/files/bz-review-report.cron
8 10 * * 1 root /usr/bin/python /usr/local/bin/bz-review-report.py -r rev-com -u --email=devel@lists.fedoraproject.org &> /dev/null
