# roles/easyfix/proxy/files/sync-easyfix.cron
5,35 * * * * root /usr/bin/rsync -a --delete sundries01::gather-easyfix/* /srv/web/gather-easyfix/ >& /dev/null
