# roles/planet/files/sub-planets/security/planet-group.cron
*/20 * * * * planet-user bash /usr/local/bin/pull-run-planet-config.sh /etc/planet/security/fpbuilder.conf >>/dev/null 2>&1
10 4 * * * *  planet-user /usr/sbin/tmpwatch 360 /srv/planet/config/security/cache >>/dev/null