# roles/openqa_worker/templates/var-lib-openqa-share.mount.j2
[Unit]
Description=openQA shared assets
ConditionPathExists=/var/lib/openqa/share

[Mount]
What={{ openqa_hostname }}:/var/lib/openqa/share
Where=/var/lib/openqa/share
Type=nfs
