# files/copr/delete-forgotten-instances.cron
#!/usr/bin/bash

source /home/copr/cloud/ec2rc.sh
/home/copr/delete-forgotten-instances.pl

