# roles/epylog/files/epylog-default.cron
#!/bin/sh
# Run epylog daily.
#
# just disable this so we don't run the global one anywhere
# better to run a specific one (or ones)
#/usr/sbin/epylog --cron
