# roles/postgresql_server/files/koji-cleanup-sessions.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root@fedoraproject.org
HOME=/

0 0 * * * postgres /usr/bin/psql -q -c "DELETE FROM sessions WHERE update_time < now() - '1 day'::interval" koji
