# roles/postgresql_server/files/koji-cleanup-locks.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root@fedoraproject.org
HOME=/

*/10 * * * * postgres kill_idle_xact_92.sh
