# files/scripts/restart-broken-ipv6.cron
* * * * * root /usr/local/bin/restart-broken-ipv6
