# roles/cloudstats/files/cloud-image-stat.cron
0 0 1 * * root /usr/bin/python /usr/local/bin/cloud-image-stat.py >& /dev/null
