# files/rdiff-backup/run-rdiff-backups.cron
# run rdiff backups
MAILTO=kevin@scrye.com
00 22 * * * root /usr/local/bin/lock-wrapper run-rdiff-backups "/usr/local/bin/run-rdiff-backups"
