# roles/supybot/files/ursabot.service
[Unit]
Description=ursabot irc bot
After=network.target

[Service]
Type=simple
User=daemon
Group=daemon
ExecStart=/usr/bin/supybot /var/lib/ursabot/conf/ursabot.conf
Restart=always
WorkingDirectory=/srv/web/meetbot

[Install]
WantedBy=multi-user.target
