# roles/mirrormanager/mirrorlist_proxy/templates/mirrorlist.service.j2
[Unit]
Description=Mirrorlist Container {{ item }}

[Service]
{% if env == "staging" %}
User=mirrormanager
{% endif %}
ExecStartPre=-/usr/bin/podman stop %n
ExecStartPre=-/usr/bin/podman rm %n --force
{% if env == "staging" %}
ExecStart=/usr/bin/podman run \
            --net=host --userns=keep-id \
            -e RUST_LOG=info --rm=true --name %n \
            -v /usr/share/GeoIP:/usr/share/GeoIP \
            -v /srv/mirrorlist/data/mirrorlist{{ item }}:/var/lib/mirrormanager \
            -v /var/log/mirrormanager:/var/log/mirrormanager \
            {{ mirrorlist_container_image_stg }} \
                 --port 1808{{ item }} \
                 --listen 127.0.01 \
                 -l /var/log/mirrormanager/%n.log
{% else %}
ExecStart=/usr/bin/podman run --rm=true --name %n -v /srv/mirrorlist/data/mirrorlist{{ item }}:/var/lib/mirrormanager:z -v /var/log/mirrormanager:/var/log/mirrormanager:z -p 1808{{ item }}:80 {{ mirrorlist_container_image }} -l /var/log/mirrormanager/%n.log
{% endif %}
ExecStop=/usr/bin/podman stop %n

[Install]
WantedBy=multi-user.target
