# roles/libravatar/files/cron.daily/libravatar
#!/bin/bash

certbot renew --post-hook 'systemctl reload httpd'
