# roles/mirrormanager/crawler/files/crawler.cron
# run the crawler twice a day
# logs sent to /var/log/mirrormanager/crawler.log and crawl/* by default
0 */12 * * * root /usr/bin/mm2_crawler --threads 75 > /dev/null 2>&1
