# roles/batcave/files/geoip-download-databases.cron
0 5 4 * * root /usr/local/bin/geoip-download-databases /srv/web/infra/bigfiles/geoip/
