# roles/batcave/files/infra-tags-report.cron
* * * * 1 nobody /usr/local/bin/infra-tags-report.sh >& /dev/null
