# roles/batcave/files/logview.cron
#!/bin/bash

# Send a email with failed or changed from ansible playbook runs 
/srv/web/infra/ansible/scripts/logview -d yesterday -s CHANGED -s FAILED | mailx -s "ansible changed/failed actions" sysadmin-logs-members@fedoraproject.org
