# roles/batcave/files/pdr.cron
0 12 * * * root /usr/bin/find /srv/web/infra/pdr -ctime +7 -print0 | xargs -0 rm -f
