# roles/batcave/files/public-db-copy.cron
0 12 * * * root /usr/local/bin/public-db-copy.sh
