# roles/batcave/files/repo2json.cron
MAILTO=sysadmin-main-members@fedoraproject.org
5 * * * * apache cd /srv/web/repo/json && /usr/bin/python /srv/web/infra/ansible/roles/repo2json/files/rhel_to_json.py > /dev/null
