# roles/batcave/files/ticketkey.cron
#!/bin/bash

# Roll over the proxy SSL ticketkey
/usr/bin/ansible-playbook /srv/web/infra/ansible/playbooks/update_ticketkey.yml | mailx -s "TickeyKey rollover" puiterwijk@fedoraproject.org
