# roles/batcave/files/ticketkey.cron
#!/bin/bash
# Roll over the proxy SSL ticketkey
source /root/sshagent >>/dev/null
export ANSIBLE_HOST_KEY_CHECKING=False
export HOME=/root/
#export ANSIBLE_SSH_PIPELINING=False
/usr/bin/ansible-playbook /srv/web/infra/ansible/playbooks/update_ticketkey.yml >/dev/null 2>&1
