# roles/mdapi/templates/mdapi.cron
# Fetches repo metadata for mdapi
#
*/20 * * * * root /usr/bin/mdapi-get_repo_md /etc/mdapi/mdapi.cfg |& grep -vi warning > /dev/null

