# roles/taskotron/imagefactory/files/imagefactory-kill-zombie.service
[Unit]
Description=Kill hanging imagefactory build processes

[Service]
Type=oneshot
ExecStart=/usr/local/bin/imagefactory-kill-zombie
TimeoutStartSec=300
