# roles/taskotron/imagefactory/files/imagefactoryd.service
# Workaround for https://github.com/redhat-imaging/imagefactory/issues/417
[Unit]
Requires=libvirtd.service
After=libvirtd.service

[Service]
Type=exec
ExecStart=/usr/bin/imagefactoryd --foreground --debug
#ExecStop=/usr/bin/killall imagefactoryd
PIDFile=/run/imagefactoryd.pid

[Install]
WantedBy=multi-user.target
