# roles/membership-map/proxy/files/sync-membership-map.cron
10,40 * * * * root /usr/bin/rsync -a --delete bapp02::membership-map/* /srv/web/membership-map/ >& /dev/null
