# roles/releng/files/ftbfs.cron
#FTBFS Weekly Reminder
* * 0 0 * ftbfs /usr/local/bin/ftbfs.sh
