# roles/releng/files/ftbfs.cron
#FTBFS Weekly Reminder
* * 0 0 * ftbfs TMPDIR=`mktemp -d /tmp/ftbfs_reminder.XXXXXX` && pushd $TMPDIR && git clone https://pagure.io/releng.git && pushd releng/scripts && ./ftbfs_weekly_reminder.py && popd && popd && popd &&rm -rf $TMPDIR
