# roles/openqa/dispatcher/files/restart-consumers.cron
#!/bin/sh

# restart all fedora messaging consumer services in case they've got stuck:
# https://github.com/fedora-infra/fedora-messaging/issues/208

/usr/bin/systemctl try-restart fm-consumer@*
