# roles/openqa/dispatcher/files/schedule-live-respins.cron
#!/bin/sh

# schedule openQA jobs for the current 'live-respins' compose; if it's
# already been tested, this will not create any new jobs

/usr/bin/fedora-openqa-schedule compose https://dl.fedoraproject.org/pub/alt/live-respins/ > /dev/null 2>&1 || /usr/bin/true
