# roles/review-stats/build/files/review-stats.cron
MAILTO=tibbs@fedoraproject.org
0,30 * * * * apache /usr/local/bin/review-stats.py -c /usr/local/share/review-stats/review-stats.cfg -t /usr/local/share/review-stats/templates -d /srv/web/review-stats/
