# roles/mirrormanager/mirrorlist2/files/mirrorlist-server.service
.include /usr/lib/systemd/system/mirrorlist-server.service
[Service]
# systemd needs an empty ExecStart= to be able to overwrite ExecStart=
ExecStart=
ExecStart=/usr/bin/python2 /usr/share/mirrormanager2/mirrorlist_server.py -l /var/log/mirrormanager/mirrorlist.log
