# roles/fedora_owner_change/files/fedora-owner-change.cron
0 10 * * 1 root /usr/bin/python /usr/local/bin/fedora-owner-change.py
