# roles/review-stats/proxy/files/sync-review-stats.cron
3,33 * * * * root /usr/bin/rsync -a --delete bapp02::review-stats/* /srv/web/review-stats/ >& /dev/null
