# roles/pdc/backend/files/pdc-audit.cron
# run once a day and report
MAILTO=releng-cron@lists.fedoraproject.org
0 18 * * * root /usr/local/bin/lock-wrapper /usr/bin/pdc-updater-audit
