# roles/mirrormanager/backend/files/backend.cron
MAILTO=root

# refresh the mirrorlist cache at the top of the hour
55 * * * * root /usr/bin/mm2_update-mirrorlist-server

# update master directory list
# logs sent to /var/log/mirrormanager/umdl.log by default
0 */2 * * * root /usr/local/bin/lock-wrapper umdl /usr/bin/mm2_update-master-directory-list

# Sync netblocks list once a day
30 0 * * * root cd /usr/share/mirrormanager2 && /usr/bin/mm2_get_global_netblocks /var/lib/mirrormanager/global_netblocks.txt
0 1 * * * root cd /usr/share/mirrormanager2 && /usr/bin/mm2_get_internet2_netblocks /var/lib/mirrormanager/internet2_netblocks.txt
30 1 * * * root /usr/bin/mm2_update-EC2-netblocks
