# roles/fedora-web/build/files/sync-atomic-htaccess.cron
MAILTO=web-members@fedoraproject.org
40 * * * * apache /usr/local/bin/lock-wrapper sync-atomic-htaccess /usr/local/bin/sync-atomic-htaccess.sh
