# roles/fedora-web/build/files/syncStatic.cron
MAILTO=web-members@fedoraproject.org
0 * * * * apache /usr/local/bin/lock-wrapper syncStatic /usr/local/bin/syncStatic
