# roles/torrent/files/torrent-web-generate.cron
*/20 * * * * torrent /srv/torrent/torrent-generator/torrent-generator
*/20 * * * * torrent /srv/torrent/spins-generator/torrent-generator

