# roles/torrent/files/torrent-hash.cron
*/20 * * * * torrent /usr/local/bin/torrent-hashes.py -o /srv/torrent/btholding/whitelist /srv/torrent/btholding/*.torrent
2 * * * * opentracker /usr/local/bin/pull_opentracker_data.sh 



