# roles/easyfix/files/easyfix.cron
# Sends the reminder about coming meetings
#
# The frequency set here should be consistent with the one set in the
# configuration file.
#
*/30 * * * * apache cd /srv/web/easyfix && /usr/bin/python /usr/local/bin/gather_easyfix.py >& /dev/null
