# roles/copr/frontend-cloud/files/systemd/system/fm-consumer@copr_messaging.service.d/50-copr.conf
[Service]
Environment=PYTHONPATH=/usr/share/copr/coprs_frontend
