# roles/blockerbugs/files/blockerbugs.cron
MAILTO="sysadmin-qa-members@fedoraproject.org"
*/30 * * * * blockerbugs blockerbugs sync > /dev/null
