# roles/blockerbugs/files/blockerbugs.cron
*/30 * * * * blockerbugs blockerbugs sync
