# roles/pkgdb2/templates/pkgdb-sync-bugzilla.cron
# Synchronize pkgdb with bugzilla
#
0 */2 * * * root PKGDB2_CONFIG=/etc/pkgdb2/pkgdb2.cfg /usr/local/bin/lock-wrapper pkgdb-sync-bugzilla "/usr/bin/pkgdb-sync-bugzilla"

