# roles/pkgdb2/templates/pkgdb-sync-bugzilla.cron
# Synchronize pkgdb with bugzilla
#
*/10 * * * * root PKGDB2_CONFIG=/etc/pkgdb2/pkgdb2.cfg /usr/bin/pkgdb-sync-bugzilla

