# roles/pkgdb2/templates/pkgdb_update_package_info.cron
# Cron updating weekly the pkgdb database to adjust summary and description of
# each package.
#
MAILTO=root

0 10 * * 5 root PKGDB2_CONFIG=/etc/pkgdb2/pkgdb2.cfg /usr/bin/update_package_info.py

