# roles/mirrormanager/frontend2/templates/frontend.cron
MAILTO=root

# sync crawler logs every hour
55 * * * * mirrormanager /usr/local/bin/sync-crawler-logs.sh
