# roles/taskotron/buildslave-configure/templates/buildslave.service.j2
[Unit]
Description=Buildslave for taskbot
After=network.target

[Service]
Type=forking
# disabled because of https://pagure.io/taskotron/issue/236
#PIDFile=/home/buildslave/slave/twistd.pid
ExecStart=/usr/bin/buildbot-worker start /home/buildslave/slave/
ExecStop=/usr/bin/buildbot-worker stop /home/buildslave/slave/
User=buildslave
Group=buildslave

[Install]
WantedBy=multi-user.target
