# roles/koji_hub/files/updatecrl.cron
0 * * * * root /usr/local/bin/updatecrl.sh &>/dev/null
