# roles/web-data-analysis/files/condense-mirrorlogs.cron
0 06  * * * root /usr/local/bin/condense-mirrorlogs.sh > /dev/null
