# roles/web-data-analysis/files/condense-hotspot.cron
0 07  * * * root /usr/local/bin/condense-hotspot.sh
