# roles/web-data-analysis/files/run-daily-awstats.cron
0 08  * * * root /usr/local/bin/run-daily-awstats.sh > /dev/null
