# roles/fedimg/files/cron/kill_ec2_nodes.cron
MAILTO=sysadmin-fedimg-members@fedoraproject.org
0 * * * * fedmsg /usr/share/fedimg/cronjobs/kill_ec2_nodes.py 2>&1 | grep -v 'Cannot find qpid python'
