# roles/fedimg/files/cron/kill_ec2_nodes.cron
MAILTO=sysadmin-fedimg-members@fedoraproject.org
0 * * * * fedmsg /usr/local/bin/kill_ec2_nodes.py
