# roles/repo2json/files/repo2json.cron
45 * * * * root cd /srv/web/repo/json && /usr/local/bin/repo2json
