# roles/easyfix/gather/templates/easyfix.cron
# Sends the reminder about coming meetings
#
# The frequency set here should be consistent with the one set in the
# configuration file.
#
{% if env == 'staging' %}
*/30 * * * * apache cd /srv/web/easyfix && /usr/bin/python /usr/local/bin/gather_easyfix.py --fedmenu-url="https://apps.stg.fedoraproject.org/fedmenu" --fedmenu-data-url="https://apps.stg.fedoraproject.org/js/data.js" >& /dev/null
{% else %}
*/30 * * * * apache cd /srv/web/easyfix && /usr/bin/python /usr/local/bin/gather_easyfix.py --fedmenu-url="https://apps.fedoraproject.org/fedmenu" --fedmenu-data-url="https://apps.fedoraproject.org/js/data.js" >& /dev/null
{% endif %}
