# roles/bodhi2/backend/files/fedwatch.service
[Service]
User=fedmsg
Group=fedmsg
Restart=on-failure
