# roles/bodhi2/backend/files/fm-consumer@.service
.include /usr/lib/systemd/system/fm-consumer@.service

[Service]
User = apache
Group = apache
