# roles/bodhi2/backend/files/new-updates-sync.cron
MAILTO=releng-cron@lists.fedoraproject.org
* * * * * ftpsync /usr/local/bin/lock-wrapper updates-push '/usr/local/bin/new-updates-sync' --silent

