# roles/distgit/files/git@.service
[Unit]
Description=Git Repositories Server Daemon
Documentation=man:git-daemon(1)

[Service]
User=nobody
ExecStart=-/usr/libexec/git-core/git-daemon --base-path=/srv/git/repositories --export-all --syslog --inetd --verbose
StandardInput=socket
