# roles/distgit/files/pkgdb_sync_git_branches.cron
00 45 * * * root /usr/local/bin/pkgdb_sync_git_branches.py
