# roles/robosignatory/files/fm-consumer@.service
[Service]
User = robosignatory
Group = robosignatory
