# roles/grobisplitter/files/rhel8-split.cron
03 09  * * * root /usr/local/bin/rhel8-split.sh
