# roles/grobisplitter/files/rhel8-split.cron
00 07  * * * root /usr/local/bin/rhel8-split.sh
