# roles/fedocal/templates/fedocal-reminder.cron
# Sends the reminder about coming meetings
#
# The frequency set here should be consistent with the one set in the
# configuration file.
#
*/30 * * * * root FEDOCAL_CONFIG=/etc/fedocal/fedocal.cfg /usr/bin/fedocal_cron.py
