# roles/postgresql_server_bdr/files/fasdb-cleanup-sessions.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root@fedoraproject.org
HOME=/

*/10 * * * * postgres /usr/local/bin/fasdb-cleanup-sessions
