# roles/grokmirror_mirror/templates/grokfsck.cron
00 02 * * * grokmirror /usr/bin/grok-fsck -c {{ grokmirror_topdir }}/fsck.conf
