# roles/people/files/check-broken-planet.cron
0 0 * * 0 /usr/local/bin/check-broken-planet.py
