# roles/people/files/del-freeipa-prcl.cron
#!/bin/bash
find /srv/groups/freeipa/prci/jobs/ -maxdepth 1 -type d -ctime +60 -exec rm -rf {} \;
