# roles/people/files/grab-daily-logs.cron
#!/bin/bash
# openstack/rdo
bash /usr/local/bin/grab-daily-logs.sh /repos/openstack/ rdo
# aarch64 image
bash /usr/local/bin/grab-daily-logs.sh F21-20140407-foundation-v8.tar.xz aarch64-stats
