# roles/people/files/grab-daily-logs.cron
#!/bin/bash
# openstack/rdo
/usr/local/bin/grab-daily-logs.sh /repos/openstack/ rdo >& /dev/null
# aarch64 image
/usr/local/bin/grab-daily-logs.sh F21-20140407-foundation-v8.tar.xz aarch64-stats >& /dev/null
