# roles/people/files/make-people-page.cron
52 * * * * apache /usr/local/bin/make-people-page.sh
