# roles/whatcanidoforfedora/build/files/sync-wcidff.cron
MAILTO=web-members@fedoraproject.org
15 * * * * apache /usr/local/bin/lock-wrapper sync-wcidff /usr/local/bin/sync-wcidff.sh
