# roles/anitya/backend/files/anitya.cron
# Checks bi-daily for new versions
#
10 */12 * * * root time ANITYA_WEB_CONFIG=/etc/anitya/anitya.cfg /usr/local/bin/lock-wrapper anitya /usr/share/anitya/anitya_cron.py
