# roles/awstats/files/combineHttp.cron
0 06  * * * root /usr/local/bin/combineHttpLogs.sh
