# roles/badges-backend/files/cron/award-badges-dev-badge.cron
40 2 * * 4 fedmsg /usr/share/badges/cronjobs/award-badges-dev-badge >> /var/log/fedmsg/cron-award-badges-dev-badge.log 2>&1
