# roles/badges-backend/files/cron/award-libravatar-badge.cron
40 * * * * fedmsg /usr/share/badges/cronjobs/award-libravatar-badge >> /var/log/fedmsg/cron-award-libravatar-badge.log 2>&1
