# roles/badges-backend/files/cron/award-libravatar-badge.cron
58 * * * * fedmsg /usr/share/badges/cronjobs/award-libravatar-badge > /var/log/fedmsg/award-libravatar-badge-cron.log 2>&1
