# roles/badges-backend/files/cron/award-lifecycle-badges.cron
40 2 * * 3 fedmsg /usr/share/badges/cronjobs/award-lifecycle-badges >> /var/log/fedmsg/cron-award-lifecycle-badges.log 2>&1
