# roles/badges-backend/files/cron/award-libravatar-badge.cron
*/55 * * * * fedmsg /usr/share/badges/cronjobs/award-libravatar-badge
