# roles/badges-backend/files/cron/award-oldschool-badges.cron
25 * * * * fedmsg /usr/share/badges/cronjobs/award-oldschool-badges > /var/log/fedmsg/award-oldschool-badges-cron.log 2>&1
