# roles/mirrormanager/mirrorlist_proxy/files/mirrorlist1.service
[Unit]
Description=Mirrorlist Container 1
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
ExecStart=/usr/bin/docker run --rm --name %n -v /srv/mirrorlist/data:/var/lib/mirrormanager:z -v /var/log/mirrormanager1:/var/log/mirrormanager:z -p 18081:80 candidate-registry.fedoraproject.org/f25/mirrormanager2-mirrorlist:f25-docker-candidate-20170206221219
ExecStop=/usr/bin/docker stop %n

[Install]
WantedBy=multi-user.target
