# roles/developer/build/files/developer-rss-update.cron
MAILTO=web-members@fedoraproject.org
45 * * * * apache /usr/local/bin/lock-wrapper developer-rss-update cd /srv/web/developer.fedoraproject.org/; python rss.py /srv/web/developer.fedoraproject.org/index.html
