# roles/developer/build/files/developer-rss-update.cron
MAILTO=web-members@fedoraproject.org
45 * * * * apache /srv/web/developer.fedoraproject.org/rss.py /srv/web/developer.fedoraproject.org/index.html
