# roles/developer/build/files/syncDeveloper.cron
MAILTO=web-members@fedoraproject.org
25 * * * * apache /usr/local/bin/lock-wrapper syncDeveloper /usr/local/bin/syncDeveloper
