# files/osbs/docker.custom.service
# Ansible managed

[Unit]
Wants=iptables.service
After=iptables.service

[Service]
ExecStartPost=/usr/local/bin/fix-docker-iptables