# files/osbs/docker.firewall.service
[Service]
ExecStartPost=/usr/local/bin/fix-docker-iptables
