# roles/koji_reminder/files/koji-cert-reminder.cron
0 0 * * 1 fedmsg /usr/local/bin/koji-cert-reminder.py 2>&1 | ts >> /var/log/koji-cert-reminder.log
