# roles/koji_reminder/files/koji-cert-reminder.cron
0 0 * * 1 fedmsg /usr/local/bin/koji-cert-reminder.py | ts >> /var/log/koji-cert-reminder.log 2>&1
